# src/tizenid/tizenid.service.in
[Unit]
Description=Making Tizen ID

[Service]
Type=oneshot
ExecStart=/usr/bin/tizen_id
ExecStartPost=/usr/bin/chsmack -a "System::Shared" @TIZEN_ID_PATH@

[Install]
WantedBy=multi-user.target
